# services/LogSphere-email.service
[Unit]
Description=Send Daily LogSphere Logs

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dailyMail.py